# chapter03/configs/alertmanager/alertdump.service
# vim: set ft=systemd :

[Unit]
Description=Alertdump logs the json payloads from Alertmanager.
After=network.target

[Service]
# Requires write permission to /vagrant/cache
User=vagrant
ExecStart=/usr/bin/alertdump
Restart=on-failure

[Install]
WantedBy=multi-user.target
